# qubes-rpc/qubes-input-sender-tablet@.service
[Unit]
Description=Qubes input proxy sender (tablet/touchscreen)
After=qubes-qrexec-agent.service

[Service]
Environment=TARGET_DOMAIN=dom0
EnvironmentFile=-/etc/qubes/input-proxy-target
ExecStart=/usr/bin/qubes-input-sender qubes.InputTablet /dev/input/%i "$TARGET_DOMAIN"
